# crontab
*/1 * * * * echo "Hello2" >> /var/log/script.log
